# torwuf/etc/cron.d/torwuf-service
#
# Regular cron jobs for the torwuf-service package
#
0 22 * * *	www-torwuf	/opt/torwuf/sbin/torwuf-service-backup
