# torwuf/etc/init/torwuf-service.conf
start on runlevel [2345]
stop on runlevel [!2345]

exec /opt/torwuf/sbin/torwuf-service

